# image/pi-gen/stage-airlock/01-airlockd/files/etc/systemd/system/airlockd.service
[Unit]
Description=Airlock daemon
Documentation=https://github.com/emdzej/airlock
After=network-online.target smbd.service avahi-daemon.service
Wants=network-online.target smbd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/airlockd
Restart=on-failure
RestartSec=3
User=root
StandardOutput=journal
StandardError=journal

# --- systemd sandboxing (seccomp + prctl only) ---
# We must NOT put airlockd in a private mount namespace: mounts made by
# the daemon there wouldn't propagate to smbd (which runs in the host
# namespace) and SMB shares would appear empty. That rules out
# ProtectSystem=, ProtectHome=, PrivateTmp=, ReadWritePaths=,
# ProtectKernelTunables=, ProtectKernelModules=, ProtectControlGroups=,
# and MountFlags= (all of which either create a mount ns or rely on one).
#
# What remains — cheap defense in depth via seccomp and prctl:
NoNewPrivileges=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

[Install]
WantedBy=multi-user.target
