# image/pi-gen/stage-airlock/01-airlockd/files/etc/systemd/system/airlockd.service
[Unit]
Description=Airlock daemon
Documentation=https://github.com/emdzej/airlock
After=network-online.target smbd.service avahi-daemon.service
Wants=network-online.target smbd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/airlockd
Restart=on-failure
RestartSec=3
User=root
StandardOutput=journal
StandardError=journal

# --- systemd sandboxing ---
# The daemon runs as root because it needs CAP_SYS_ADMIN to mount/unmount
# block devices and spawn mkfs / parted / smbcontrol. These directives
# constrain what "root" can reach: read-only view of /usr and /etc except
# for the specific paths we regenerate, no /home, private /tmp, no kernel
# tunables / modules / cgroups, no namespace or personality tricks.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/etc/samba/smb.conf.d /mnt/airlock /var/log/samba
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
